# systemd/jellyfin-stack-rclone@.service.d/overlay.conf
[Service]
Environment=RCLONE_MEDIA_ALLOW_OVERLAY=1
SuccessExitStatus=143 SIGTERM
