# systemd/jellyfin-stack.service.d/boot-profiles.conf
[Service]
SuccessExitStatus=143 SIGTERM
Restart=on-failure
RestartSec=20
ExecStart=
ExecStart=%h/.local/bin/jellyfin-stack-boot.sh --env-file %h/jellyfin-stack/.env
