# etc/systemd/system/openlinkhub.service.d/override.conf
[Unit]
# Restart even if stopped cleanly or by another agent session
StartLimitIntervalSec=0
After=multi-user.target

[Service]
# Wait for USB devices to enumerate after boot
ExecStartPre=/bin/sleep 10
# Apply Performance profile after startup
ExecStartPost=/bin/bash -c 'sleep 20 && curl -s -X POST http://127.0.0.1:27003/api/speed -d '"'"'{"deviceId":"45C1646A377CE85BB48901FD77EB892A","channelId":16,"profile":"Performance"}'"'"' && curl -s -X POST http://127.0.0.1:27003/api/speed -d '"'"'{"deviceId":"45C1646A377CE85BB48901FD77EB892A","channelId":15,"profile":"Performance"}'"'"''
# Always restart, no rate limiting
Restart=always
RestartSec=5
