# etc/systemd/system/jellyfin.service.d/10-optional-env.conf
[Service]
Environment="JELLYFIN_SERVICE_OPT="
Environment="JELLYFIN_NOWEBAPP_OPT="
