# etc/systemd/system/jellyfin.service.d/20-webroot-path.conf
[Service]
ExecStartPre=/usr/bin/ln -sfn /usr/share/jellyfin/web /var/lib/jellyfin/wwwroot
