# systemd/bar-cve.timer
[Unit]
Description=Refresh CVE-advisory cache hourly

[Timer]
OnBootSec=2min
OnUnitActiveSec=1h
AccuracySec=2min

[Install]
WantedBy=timers.target
