# systemd/bar-prs.timer
[Unit]
Description=Refresh GitHub PRs cache every 5 minutes

[Timer]
OnBootSec=60s
OnUnitActiveSec=5min
AccuracySec=30s

[Install]
WantedBy=timers.target
