# systemd/dotfiles-keybind-ticker.service
[Unit]
Description=dotfiles keybind ticker (pixel-smooth scrolling bar)
PartOf=graphical-session.target
After=graphical-session.target ironbar.service

[Service]
Type=simple
ExecCondition=/usr/bin/sh -c 'test -n "$WAYLAND_DISPLAY" && test -S "$XDG_RUNTIME_DIR/$WAYLAND_DISPLAY"'
ExecStartPre=/usr/bin/sh -c 'for i in $(seq 1 20); do hyprctl version > /dev/null 2>&1 && exit 0; sleep 0.5; done; exit 1'
# Smoke test each plugin's META + build() + Pango markup validity
# before starting; catches unescaped &/< in dynamic text and broken
# plugins at boot time instead of on first rotation. --fast skips the
# slow-threaded streams (github/music/updates/claude-sessions) so
# startup adds ~200 ms in the common case.
ExecStartPre=/usr/bin/python3 %h/hairglasses-studio/dotfiles/scripts/ticker-smoke-test.py --fast
Environment=PATH=%h/.local/bin:%h/bin:/usr/local/bin:/usr/bin:/bin
Environment=GSK_RENDERER=gl
Environment=LD_PRELOAD=/usr/lib/libgtk4-layer-shell.so
ExecStart=/usr/bin/python3 %h/hairglasses-studio/dotfiles/scripts/keybind-ticker.py --layer
SuccessExitStatus=143
Restart=on-failure
RestartSec=2
# Crash-loop guard: after 5 starts in 60s systemd stops retrying so we
# don't mask a genuine bug behind a silent retry loop. Prior default
# (5 in 10s) was too aggressive — brief hyprland reloads tripped it.
StartLimitBurst=5
StartLimitIntervalSec=60
# Resource caps: ticker peaks at ~100M RSS after glyph-mask caching;
# 160M leaves breathing room for phosphor_trail + ghost_echo frames.
MemoryMax=160M
MemoryHigh=128M
CPUQuota=200%
# Hardening — the ticker is a pure consumer; it never writes outside
# /tmp and ~/.local/state/keybind-ticker. ProtectSystem=strict is too
# heavy because GTK needs /etc/gtk-4.0 etc; ProtectHome=read-only is
# wrong because we write state files. Stick to the low-risk flags.
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=graphical-session.target
