# systemd/dotfiles-kitty-save-session.service
[Unit]
Description=dotfiles kitty session saver
# Optional helper for explicit kitty-save-session users. The default Kitty
# launch policy still comes from kitty.conf (startup_session none).
PartOf=graphical-session.target
After=graphical-session.target
ConditionPathExists=%h/.config/kitty/extras/kitty-save-session/kitty-save-session-all.sh

[Service]
Type=oneshot
Environment=KITTY_SESSION_SOCK_PATTERN=@mykitty-*
Environment=KITTY_SESSION_SAVE_DIR=%h/.local/state/kitty/sessions
Environment=KITTY_SESSION_SAVE_OPTS=--no-copy-env
ExecStart=%h/.config/kitty/extras/kitty-save-session/kitty-save-session-all.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dotfiles-kitty-save-session

[Install]
WantedBy=graphical-session.target
