# systemd/dotfiles-notification-history.service
[Unit]
Description=dotfiles notification history listener
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
Environment=XDG_RUNTIME_DIR=%t
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=%t/bus
ExecStart=/usr/bin/env python3 %h/hairglasses-studio/dotfiles/scripts/notification-history-listener.py
Restart=always
RestartSec=2

[Install]
WantedBy=graphical-session.target
