# systemd/dotfiles-hyprland-autoname-workspaces.service
[Unit]
Description=dotfiles workspace autonamer for Hyprland
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
Environment=PATH=%h/.local/bin:%h/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=%h/hairglasses-studio/dotfiles/scripts/run-wayland-session-command.sh -- hyprland-autoname-workspaces
SuccessExitStatus=143
Restart=always
RestartSec=10

[Install]
WantedBy=graphical-session.target
