# systemd/dotfiles-hypr-dock.service
[Unit]
Description=dotfiles hypr-dock bottom launcher
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecCondition=/usr/bin/sh -c 'test -n "$WAYLAND_DISPLAY" && test -S "$XDG_RUNTIME_DIR/$WAYLAND_DISPLAY"'
Environment=PATH=%h/.local/bin:%h/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=%h/hairglasses-studio/dotfiles/scripts/run-hypr-dock.sh -- -config %h/.config/hypr-dock/hypr-dock.conf
SuccessExitStatus=143
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
