# systemd/dotfiles-kitty-save-session.timer
[Unit]
Description=dotfiles kitty session save timer

[Timer]
OnActiveSec=2m
OnUnitActiveSec=10m

[Install]
WantedBy=timers.target
