# systemd/juhradial-overlay.service
[Unit]
Description=JuhRadial MX Overlay
After=graphical-session.target juhradialmx-daemon.service
PartOf=graphical-session.target

[Service]
Type=simple
ExecStartPre=/usr/bin/env sh -c 'mkdir -p "%h/.local/state/juhradial"'
ExecStart=%h/hairglasses-studio/dotfiles/scripts/run-juhradial-overlay.sh
StandardOutput=append:%h/.local/state/juhradial/overlay.log
StandardError=append:%h/.local/state/juhradial/overlay.log
SuccessExitStatus=143 SIGTERM
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
