# systemd/juhradialmx-daemon.service
[Unit]
Description=JuhRadial MX Daemon
Documentation=https://github.com/JuhLabs/juhradial-mx
After=graphical-session.target
Wants=graphical-session.target
PartOf=graphical-session.target
Requires=dbus.socket
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
ExecStart=%h/.local/bin/juhradiald
Restart=on-abnormal
RestartSec=5s
WatchdogSec=0
Environment=RUST_LOG=info
Environment=RUST_BACKTRACE=1
NoNewPrivileges=yes
MemoryMax=100M
CPUQuota=50%

[Install]
WantedBy=default.target
