# systemd/mx-battery-notify.service
[Unit]
Description=MX Master 4 battery check

[Service]
Type=oneshot
ExecStart=%h/hairglasses-studio/dotfiles/scripts/mx-battery-notify.sh
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus
