# systemd/rclone-mega.service
[Unit]
Description=Mount MEGA library with rclone
After=network-online.target
Wants=network-online.target
OnFailure=notify-failure@%n.service

[Service]
Type=notify
EnvironmentFile=%h/.config/stashapp/runtime.env
ExecStartPre=/usr/bin/env bash -lc '/bin/fusermount -uz "${MOUNT_DIR}" >/dev/null 2>&1 || /bin/fusermount3 -uz "${MOUNT_DIR}" >/dev/null 2>&1 || /usr/bin/umount -l "${MOUNT_DIR}" >/dev/null 2>&1 || true'
ExecStart=/usr/bin/rclone mount ${RCLONE_SOURCE} ${MOUNT_DIR} --config ${RCLONE_CONFIG} --vfs-cache-mode full --vfs-cache-max-size 50G --vfs-cache-max-age 72h --vfs-read-chunk-size 32M --vfs-read-chunk-streams 4 --vfs-write-back 5s --buffer-size 32M --dir-cache-time 5m --attr-timeout 10s --poll-interval 1m --transfers 4 --checkers 4 --cache-dir ${VFS_CACHE_DIR} --use-mmap --no-modtime --retries 3
ExecStop=/usr/bin/env bash -lc '/bin/fusermount -u "${MOUNT_DIR}" >/dev/null 2>&1 || /bin/fusermount -uz "${MOUNT_DIR}" >/dev/null 2>&1 || /bin/fusermount3 -u "${MOUNT_DIR}" >/dev/null 2>&1 || /bin/fusermount3 -uz "${MOUNT_DIR}" >/dev/null 2>&1 || /usr/bin/umount -l "${MOUNT_DIR}" >/dev/null 2>&1 || true'
Restart=on-failure
RestartSec=10
LimitNOFILE=32768
TimeoutStopSec=30

[Install]
WantedBy=default.target
