# systemd/rg-marathon@.service
[Unit]
Description=Ralphglasses Marathon for %i
After=default.target
StartLimitIntervalSec=600
StartLimitBurst=5
OnFailure=notify-failure@%n.service

[Service]
Type=simple
WorkingDirectory=%h/hairglasses-studio/%i
ExecStart=%h/hairglasses-studio/ralphglasses/marathon.sh -b 100 -d 12 -c 80
Restart=on-failure
RestartSec=30
SyslogIdentifier=rg-marathon-%i
MemoryMax=4G

[Install]
WantedBy=default.target
