# systemd/shader-rotate.service
[Unit]
Description=Advance kitty shader rotation
ConditionPathExists=%t/hypr

[Service]
Type=oneshot
ExecStart=/usr/bin/env bash %h/hairglasses-studio/dotfiles/scripts/kitty-shader-playlist.sh next ambient
Environment=PATH=%h/.local/bin:/usr/bin:/usr/local/bin:/bin
