# systemd/shader-rotate.service
[Unit]
Description=Advance Ghostty shader rotation

[Service]
Type=oneshot
ExecStart=%h/hairglasses-studio/dotfiles/ghostty/shaders/bin/shader-auto-rotate.sh
Environment=PATH=/usr/bin:/bin
