# systemd/shader-rotate.service
[Unit]
Description=Advance Ghostty shader rotation

[Service]
Type=oneshot
ExecStart=%h/dotfiles/ghostty/shaders/bin/shader-auto-rotate.sh
