# systemd/shader-rotate.timer
[Unit]
Description=Rotate Kitty visuals every 30 minutes

[Timer]
OnBootSec=60
OnUnitActiveSec=1800

[Install]
WantedBy=timers.target
