# systemd/shader-rotate.service
[Unit]
Description=Advance Ghostty shader rotation
ConditionPathExists=%t/hypr

[Service]
Type=oneshot
ExecStart=/usr/bin/env zsh %h/hairglasses-studio/dotfiles/ghostty/shaders/bin/shader-auto-rotate.sh
Environment=PATH=%h/.local/bin:/usr/bin:/usr/local/bin:/bin
