# systemd/shader-rotate.service
[Unit]
Description=Advance kitty visual rotation
ConditionPathExists=%t/hypr

[Service]
Type=oneshot
ExecStart=%h/.local/bin/kitty-shader-playlist next ambient
Environment=PATH=%h/.local/bin:/usr/bin:/usr/local/bin:/bin
