# systemd/stash.service
[Unit]
Description=Stash media organizer
After=network-online.target rclone-mega.service
Wants=network-online.target
Requires=rclone-mega.service
OnFailure=notify-failure@%n.service

[Service]
Type=simple
WorkingDirectory=%h/mega/linux-isos/stashapp
EnvironmentFile=%h/.config/stashapp/runtime.env
ExecStart=/usr/bin/env bash %h/mega/linux-isos/stashapp/run-stash.sh run
ExecStop=/usr/bin/env bash %h/mega/linux-isos/stashapp/run-stash.sh stop
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=default.target
