# systemd/stash-ai-server.service
[Unit]
Description=Stash AI backend
After=network-online.target stash.service
Wants=network-online.target
Requires=stash.service

[Service]
Type=simple
WorkingDirectory=%h/mega/linux-isos/stashapp
EnvironmentFile=%h/.config/stashapp/runtime.env
ExecStart=/usr/bin/env bash %h/mega/linux-isos/stashapp/tools/stash-ai-server-backend/run-stash-ai-server.sh run
ExecStop=/usr/bin/env bash %h/mega/linux-isos/stashapp/tools/stash-ai-server-backend/run-stash-ai-server.sh stop
Restart=on-failure
RestartSec=10
TimeoutStopSec=60

[Install]
WantedBy=default.target
