# mcp/hg-mcp/com.aftrs.music-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aftrs.music-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>~/aftrs-studio/hg-mcp/bin/music-sync</string>
        <string>--all</string>
    </array>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>WorkingDirectory</key>
    <string>~/aftrs-studio/hg-mcp</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>AWS_PROFILE</key>
        <string>cr8</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/music-sync.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/music-sync.error.log</string>

    <key>Nice</key>
    <integer>10</integer>

    <key>LowPriorityIO</key>
    <true/>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
